# systemd/openwheel-daemon.service.in
# SPDX-License-Identifier: GPL-3.0-or-later
# openwheel-daemon systemd user service
#
# Install system-wide: sudo cmake --install build --prefix /usr
# Then enable for your user:
#   systemctl --user enable --now openwheel-daemon
#
# The daemon requires /dev/hidraw* and /dev/uinput access.
# Ensure the udev rule is loaded first:
#   sudo cp udev/99-openwheel.rules /etc/udev/rules.d/
#   sudo udevadm control --reload && sudo udevadm trigger
#   sudo usermod -aG input $USER   # log out and back in

[Unit]
Description=OpenWheel Dial Daemon
Documentation=https://github.com/fredaime/openwheel
# Start after the graphical session is ready (session D-Bus available)
After=graphical-session.target
PartOf=graphical-session.target
# Stop retrying after 5 consecutive failures in 60 s — prevents D-Bus hammering
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/asus_wheel
# Only restart on unexpected failures, not on clean exit (e.g. device not found)
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=graphical-session.target
